# cron.d/crontab
* * * * * echo "Cron is Running in cron" >> /tmp/cron.log
* * * * * php /var/www/html/artisan schedule:run >> /dev/null 2>&1
